# root/etc/service/btsync/run
#!/bin/sh
exec 2>&1
exec /usr/local/bin/btsync --config /etc/btsync/btsync.conf --log /var/log/btsync.log --nodaemon
